# bqclock.timer
[Unit]
Description=Run bqclock every minute.

[Timer]
OnActiveSec=1min
OnUnitActiveSec=1min
Unit=bqclock.timer

[Install]
WantedBy=timers.target
