# bqclock.service
[Unit]
Description=BeQuiet Keyboard Clock Sync

[Service]
ExecStart=/usr/local/bin/bqclock
